# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=5MB
loglevel=info

[program:imap_poller]
command=python -m app.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:webui]
# Set WEBUI_ENABLED=false in docker-compose.yml to disable the web interface.
command=bash -c 'if [ "${WEBUI_ENABLED:-true}" = "true" ]; then uvicorn app.webui.server:app --host 0.0.0.0 --port 635; else echo "WebUI disabled (WEBUI_ENABLED=false)"; sleep infinity; fi'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
