# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/tmp/supervisord.pid

[program:mail2print]
command=python /app/app/main.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/mail2print.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/mail2print_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:webui]
command=python /app/app/webui/server.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/webui.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/app/logs/webui_err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
