# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/tmp/supervisord.pid

[program:mail2print]
; Als Modul starten, damit relative Imports (from app.xxx) funktionieren
command=python -m app.main
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/mail2print.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/mail2print_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:webui]
; FastAPI-App muss über uvicorn als ASGI-Server gestartet werden
command=uvicorn app.webui.server:app --host 0.0.0.0 --port 635
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/webui.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/app/logs/webui_err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
